# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:backend]
command=python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/app/backend
environment=DATABASE_URL="sqlite:////app/data/talentlens.db",UPLOAD_DIR="/app/data/uploads",CHROMA_DIR="/app/data/chroma"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:frontend]
command=npm start
directory=/app/frontend
environment=PORT="3000",HOSTNAME="0.0.0.0",BACKEND_URL="http://127.0.0.1:8000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
